# templates/default/etc/initd.erb
#!/bin/bash
#
# chkconfig: 235 10 90
# description: TeamCity startup script
#

TEAMCITY_SERVICE_NAME=<%= @tc_service_name %>
TEAMCITY_USER=<%= @tc_user %>
TEAMCITY_AGENT=<%=@tc_agent_home%>/bin/agent.sh
TEAMCITY_PIDFILE=<%=@tc_agent_home%>/logs/buildAgent.pid

. /etc/rc.d/init.d/functions

start() {
  su - $TEAMCITY_USER sh -c "$TEAMCITY_AGENT start"
}

stop() {
  su - $TEAMCITY_USER sh -c "$TEAMCITY_AGENT stop"
}

restart() {
  stop
  start
}

script_result=0

case "$1" in
  start)
    start
    ;;
  status)
    status -p $TEAMCITY_PIDFILE $TEAMCITY_SERVICE_NAME
    script_result=$?
    ;;
  stop)
    stop
    ;;
  restart)
    restart
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    script_result=1
    ;;
esac

exit $script_result

