# templates/default/etc/init.erb
#!/bin/bash
# chkconfig: 234 20 80
ps -ef | grep java | grep jetbrains.buildServer.agent.AgentMain | grep -v grep >> /dev/null 2>&1

RUNNING=$?

set -e

INITSCRIPT="$(basename "$0")"
JOB="${INITSCRIPT%.sh}"

if [ -z "$1" ]; then
        echo "Usage: $0 COMMAND" 1>&2
        exit 1
fi

COMMAND="$1"
shift

case $COMMAND in
start|stop)
    if [ "$RUNNING" -eq 0 ] && [ "$COMMAND" = "start" ]; then
        exit 0
    elif [ "$RUNNING" -ne 0 ] && [ "$COMMAND" = "stop" ]; then
        exit 0
    fi
    su -l teamcity teamcity bash -c "
    <%= @tc_agent_home %>/bin/agent.sh $COMMAND"
    ;;
restart)
    export RUNNING ;
    su -l  teamcity bash -c "
    #source \"/usr/local/rvm/scripts/initialize\" ;
    if [ \"$RUNNING\" -eq 0 ]; then
        <%= @tc_agent_home%>/bin/agent.sh stop kill;
    fi;
    <%= @tc_agent_home %>/bin/agent.sh start"
    ;;
*)
    exit 1
esac
